# cronjobs/crontab
* * * * * python /assembly_import.py
* * * * * python /biosample_import.py
* * * * * python /biosample_update.py 
* * * * * python /countries_update.py 
* * * * * python /read_update.py