# cronjobs/crontab

* * * * * python /countries_update.py 
* * * * * python /locations_update.py