# cronjobs/crontab

* * * * * python /cronjobs/assembly_import.py 