# cronjobs/crontab
* * * * * python /trigger_cronjob.py helpers unset_fields
