# cronjobs/crontab
* * * * * python /trigger_cronjob.py helpers add_lineage
* * * * * python /trigger_cronjob.py helpers handle_orphans
