# cronjobs/crontab
* * * * * python /trigger_cronjob.py import_assemblies

