# cronjobs/bulk_operations/crontab
* * * * * /usr/local/bin/python /cronjobs/assembly_import.py > /proc/1/fd/1 2>/proc/1/fd/2
